# configs/systemd/camilladsp.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/camilladsp -a 127.0.0.1 -p 1234 /etc/camilladsp/active.yml
