# configs/systemd/user/pi4audio-signal-gen.service
[Unit]
Description=RT Signal Generator for Pi Audio Workstation (D-037, D-040)
# D-040: GraphManager controls all link topology. Signal-gen runs in managed
# mode — it registers PW streams with ports but does not auto-connect.
# GM must be running before signal-gen starts so links can be established.
# WirePlumber activates adapter node ports (hardware device profiles).
After=pipewire.service wireplumber.service pi4audio-graph-manager.service
Requires=pipewire.service
Wants=wireplumber.service pi4audio-graph-manager.service

[Service]
Type=simple
ExecStart=%h/bin/pi4audio-signal-gen \
    --managed \
    --capture-target UMIK-1 \
    --channels 1 \
    --listen tcp:127.0.0.1:4001 \
    --max-level-dbfs=-20.0

Restart=on-failure
RestartSec=2

# RT scheduling: inherited from PipeWire.
# The signal generator uses RT_PROCESS streams, so PipeWire's data thread
# invokes our process callbacks at its own RT priority (SCHED_FIFO/88 on
# this system). We do NOT set CPUSchedulingPolicy here -- the binary itself
# runs at normal priority; only the PW data thread callbacks are RT.
# This matches how all PipeWire clients (Mixxx, Reaper) work: the application
# thread is SCHED_OTHER, the PW callback thread is SCHED_FIFO.

# Prevent OOM killer from targeting us during audio processing
OOMScoreAdjust=-500

# Resource limits
LimitMEMLOCK=infinity
LimitRTPRIO=88

# Security hardening — PipeWire client profile (SEC-PW-CLIENT)
# See F-035 for rationale.
NoNewPrivileges=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=default.target
